# deployment/systemd-queue-worker@.service
#==============================================================================
# Systemd Service for Queue Workers
#==============================================================================
#
# Installation:
#   sudo cp deployment/systemd-queue-worker@.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable queue-worker@default
#   sudo systemctl enable queue-worker@high-priority
#   sudo systemctl start queue-worker@default
#   sudo systemctl start queue-worker@high-priority
#
# Management:
#   sudo systemctl status queue-worker@default
#   sudo systemctl restart queue-worker@default
#   sudo systemctl stop queue-worker@default
#   sudo journalctl -u queue-worker@default -f
#
# Multiple instances:
#   systemctl start queue-worker@default
#   systemctl start queue-worker@high-priority
#   systemctl start queue-worker@emails
#
# IMPORTANT: Update paths below to match your project!
#==============================================================================

[Unit]
Description=Queue Worker %i
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/html

# %i is replaced with instance name (default, high-priority, etc.)
ExecStart=/usr/bin/php /var/www/html/worker.php %i 0 3

# Auto-restart on failure
Restart=always
RestartSec=5

# Timeout for graceful shutdown
TimeoutStopSec=3600

# Logging
StandardOutput=append:/var/www/html/storage/logs/worker-%i.log
StandardError=append:/var/www/html/storage/logs/worker-%i.log

# Resource limits (optional)
LimitNOFILE=65536
MemoryLimit=512M

# Environment variables (optional)
Environment="APP_ENV=production"

[Install]
WantedBy=multi-user.target
